# kafka/supervisor/kafka.conf
[program:kafka]
command=%(ENV_KAFKA_HOME)s/bin/kafka-server-start.sh %(ENV_KAFKA_HOME)s/config/server.properties
autostart=true
autorestart=true
