# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=%(ENV_KAFKA_HOME)s/bin/zookeeper-server-start.sh %(ENV_KAFKA_HOME)s/config/zookeeper.properties
autostart=true
autorestart=true